# database/crontab
0 * * * * find /var/unsavedtexts -mtime +59m -type f -exec rm {} \; >/var/log/unsaved-text-cleanup.log 2>&1
